# launchd/com.polymarket.position_mgr.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.polymarket.position_mgr</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/zsh</string>
        <string>-lc</string>
        <!-- Update path to match your actual repo location -->
        <string>cd /Users/damirkurecic/Documents/github/polymarket-bot &amp;&amp; venv/bin/python position_manager.py --execute &gt;&gt; /tmp/polymarket_position_mgr.log 2&gt;&amp;1</string>
    </array>

    <!-- Run every 5 minutes -->
    <key>StartInterval</key>
    <integer>300</integer>

    <key>RunAtLoad</key>
    <false/>

    <key>StandardOutPath</key>
    <string>/tmp/polymarket_position_mgr.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/polymarket_position_mgr_err.log</string>
</dict>
</plist>
